# tools/labs/templates/device_model/sol/parallel/parallel.rules

SUBSYSTEM=="pnp", ATTRS{id}=="PNP0400", RUN+="/sbin/insmod /root/pso/lab12/pnp/mydev/mydev.ko"

